# lxqt-policykit/autostart/lxqt-policykit-agent.desktop.in
[Desktop Entry]
Type=Application
Name=PolicyKit Handler
TryExec=lxqt-policykit-agent
Exec=lxqt-policykit-agent
OnlyShowIn=LXQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
