# lxqt-panel/autostart/lxqt-panel.desktop.in
[Desktop Entry]
Type=Application
Name=Panel
TryExec=lxqt-panel
Exec=lxqt-panel
OnlyShowIn=LXQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
