# lxqt-globalkeys/autostart/lxqt-globalkeyshortcuts.desktop.in
[Desktop Entry]
Type=Application
Name=Global Keyboard Shortcuts
TryExec=lxqt-globalkeysd
Exec=lxqt-globalkeysd
OnlyShowIn=LXQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
